# packaging/keyway-visualizer.service
[Unit]
Description=Keyway Visualizer
After=graphical-session.target

[Service]
ExecStart=keyway-visualizer
Restart=on-failure
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
